# opt/retropie/configs/all/CRT/bin/BackGroundMusic/BackGroundMusic.service
[Unit]
Description=This service play background music with emulationstionation
DefaultDependencies=no

[Service]
ExecStart=/usr/bin/python -u /opt/retropie/configs/all/CRT/bin/BackGroundMusic/BGM.py
WorkingDirectory=/opt/retropie/configs/all/CRT/bin/BackGroundMusic
StandardOutput=inherit
StandardError=inherit
Restart=no
User=pi

[Install]
WantedBy=sysinit.target
