# opt/retropie/configs/all/CRT/bin/ScreenUtilityFiles/bin/module_display/oled.service
[Unit]
Description=This service handles piCRT OLED display
DefaultDependencies=no

[Service]
ExecStart=/usr/bin/python3 -u oled.py
WorkingDirectory=/opt/retropie/configs/all/CRT/bin/ScreenUtilityFiles/bin/module_display
StandardOutput=inherit
StandardError=inherit
Restart=no
User=pi

[Install]
WantedBy=sysinit.target