# opt/retropie/configs/all/CRT/bin/ScreenUtilityFiles/bin/module_cable/rgb_cable.service
[Unit]
Description=This service will enable jamma controls software and other boot video settings

[Service]
ExecStart=/usr/bin/python3 -u rgb_cable.py
WorkingDirectory=/opt/retropie/configs/all/CRT/bin/ScreenUtilityFiles/bin/module_cable
StandardOutput=null
StandardError=null
Restart=no
User=pi

[Install]
WantedBy=basic.target